# apps/worker/sentinel-worker.service
[Unit]
Description=Sentinel Worker Service
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/sentinel/worker
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=append:/path/to/sentinel/worker/logs/worker.log
StandardError=append:/path/to/sentinel/worker/logs/worker.log
Environment="NODE_ENV=production"
EnvironmentFile=/path/to/sentinel/worker/.env.local

[Install]
WantedBy=multi-user.target
